# tools/archived-harvester/systemd/harvester-daily-open.service
[Unit]
Description=Trusty Squire skill harvester — open daily summary issue

[Service]
Type=oneshot
WorkingDirectory=/home/lunchbox/trusty-squire/tools/harvester
ExecStart=/usr/bin/node /home/lunchbox/trusty-squire/tools/harvester/daily.mjs --create
Environment=GH_REPO=Trusty-Squire/trusty-squire
StandardOutput=journal
StandardError=journal
