# tools/systemd/trusty-squire-heal.service
# The registry verify pass — a systemd USER service on the operator's box.
# Oneshot: the .timer triggers ONE verify pass per fire. The pass asks codex
# (holding the @next Trusty Squire MCP) to reproduce each queued skill's signup;
# the registry then applies the mechanical rule (one success promotes pending-
# review → active, the 3rd consecutive real failure demotes). See
# docs/DESIGN-housekeeper-codex-verify.md.
#
# Install: see tools/systemd/README.md. Requires `loginctl enable-linger`.
[Unit]
Description=Trusty Squire registry verify pass (codex-driven promote/demote)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Runs from a SOURCE checkout (the package is operator-only, never published).
# Keep this checkout on `staging` so the run dogfoods the latest code. Adjust
# WorkingDirectory if your checkout isn't ~/proj-ts. The package builds to
# apps/housekeeper/dist (pnpm -F @trusty-squire/housekeeper build).
WorkingDirectory=%h/proj-ts
# Operator env: REGISTRY_ADMIN_BEARER (required), TRUSTY_SQUIRE_REGISTRY_URL,
# and the optional HOUSEKEEPER_* knobs (MAX_SKILLS_PER_RUN, CODEX_CMD,
# PROMOTE_MIN_SUCCESSES). NO machine token / LLM tier / inbox creds — the verify
# pass has no discover phase and no LLM proxy; codex is the only intelligence.
EnvironmentFile=%h/.config/trusty-squire/harvester.env
# Prereq gate — fail fast (and loudly) if the codex CLI isn't installed/on PATH
# or the admin bearer is missing, instead of silently no-op'ing the whole pass.
ExecStartPre=/usr/bin/env node apps/housekeeper/dist/bin.js --check
# The verify pass. --once = a single pass (the timer provides the cadence).
ExecStart=/usr/bin/env node apps/housekeeper/dist/bin.js --once
# codex driving a full browser signup can take minutes; with the default
# 20-skill cap that's up to ~3h. Give room, but don't let a hung run wedge it.
TimeoutStartSec=4h
# Don't thrash on a transient failure — the next timer fire retries.
Restart=no
