# tools/systemd/trusty-squire-heal.service
# T9 — the closed-loop self-healing pass (DESIGN-closed-loop-remediation.md).
# A systemd USER service (not system-wide): runs on the operator's headless
# dev box as the operator, with their bot Chrome profile + env. Oneshot —
# the .timer triggers one heal pass (verify→discover→digest) per fire.
#
# Install: see tools/systemd/README.md. Requires `loginctl enable-linger`
# so the timer survives logout.
[Unit]
Description=Trusty Squire closed-loop heal pass (verify rot, re-skill, digest)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# The housekeeper runs from a SOURCE checkout (it's excluded from the npm
# tarball). Adjust WorkingDirectory if your checkout isn't ~/proj-ts.
#
# C4 — DOGFOOD THE RC: keep this checkout on the `staging` branch so each daily
# run executes the latest release candidate (the `next` dist-tag). That makes
# the run the live test of yesterday's RC, and the OF#2 it reports the per-RC
# promote signal on the dashboard (see docs/DESIGN-autonomous-output-loop.md).
WorkingDirectory=%h/proj-ts
# Operator env: REGISTRY_ADMIN_BEARER (required), TRUSTY_SQUIRE_REGISTRY_URL,
# UNIVERSAL_BOT_LLM_TIER=free, optional TELEGRAM_BOT_TOKEN for the digest.
# The machine token + account id (discover phase) do NOT need to live in
# harvester.env — the CLI auto-reads them from ~/.config/trusty-squire/
# session.json (an explicit env value still wins if you set one).
#
# Recommended in harvester.env: UNIVERSAL_BOT_MAX_LLM_CALLS=25. Discovery hits
# hard multi-step onboarding wizards; the default 15-call breaker guillotines
# legitimately-progressing wizards mid-flight (measured 2026-06-11: anyscale +
# kinde cut off at 15 while still advancing). Housekeeper-only — end-user
# installs keep 15. The loop-detector still bails true planning loops early.
EnvironmentFile=%h/.config/trusty-squire/harvester.env
# The daily autonomous engine: verify (keep the registry honest) → discover
# across the curated ~300-service candidate queue (maximize skills) → one digest
# carrying both objective functions (skills-in-registry + discovery success
# rate). --from sources the curated sweep; --limit caps how many services one
# pass attempts (kept well under what fits the 8h timeout at ~6min/signup).
# skip-active drops already-skilled services at runtime, so the freed budget
# always goes to net-new coverage. Drop --from to fall back to on-demand
# telemetry-candidate discovery instead of the curated sweep.
ExecStart=/usr/bin/env node apps/mcp/dist/bin.js housekeeper --mode=heal --once --telegram --from=tools/discovery-candidates.yaml --limit=80
# C2 — the OUTPUT side of the loop. After the run measures, the holistic
# fix-agent reads the failure batch from the capture dir, clusters by root
# cause, proposes a generalizing fix per cluster, reconciles it against the
# planner eval gate, and commits surviving fixes to staging (the `next`/RC
# channel). Best-effort (leading `-`): a fix-pass failure never fails the timer.
#   • Needs a coding CLI: TRUSTY_SQUIRE_FIX_AGENT_CLI (default "claude -p").
#   • PUSH IS OPT-IN: without TRUSTY_SQUIRE_FIX_AGENT_PUSH=1 it commits RCs to
#     the local staging branch but does NOT push — nothing reaches npm until
#     the operator pushes (or sets the env). next→latest stays a human promote.
ExecStartPost=-/usr/bin/env node apps/mcp/dist/bin.js housekeeper --mode=fix
# A heal pass drives ~100 real browser signups in the discover phase; give it
# plenty of room but don't let a hung run wedge the timer.
TimeoutStartSec=8h
# Don't thrash on a transient failure — the next timer fire retries.
Restart=no
