# tools/systemd/trusty-squire-heal.timer
# Fires the daily registry verify pass — one codex-driven pass per day over the
# verifier queue (promote pending-review, demote stale).
# Persistent=true → if the box was asleep/off at the scheduled time, it runs
# on next boot rather than skipping a cycle.
[Unit]
Description=Run the Trusty Squire registry verify pass daily

[Timer]
OnBootSec=10min
OnUnitActiveSec=24h
Persistent=true
# Spread the start a little so it never lines up exactly with other timers.
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
