# tools/systemd/trusty-squire-heal.timer
# Fires the daily autonomous engine — one heal pass per day over the curated
# ~100-service queue, emitting one digest with both objective functions.
# Persistent=true → if the box was asleep/off at the scheduled time, it runs
# on next boot rather than skipping a cycle. (24h keeps the dashboard's 26h
# DOWN-staleness threshold valid.)
[Unit]
Description=Run the Trusty Squire autonomous heal pass daily (~100 services)

[Timer]
OnBootSec=10min
OnUnitActiveSec=24h
Persistent=true
# Spread the start a little so it never lines up exactly with other timers.
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
