# tools/systemd/trusty-squire-heal.timer
# Fires the daily registry verify pass — one codex-driven pass per day over the
# verifier queue (promote pending-review, demote stale).
#
# OnCalendar (not OnBootSec/OnUnitActiveSec): a calendar schedule always computes
# a next-elapse, and — unlike a monotonic timer — Persistent=true actually works
# with it, so a missed run (box asleep/off at 04:00) fires on the next wake
# instead of being skipped. The old monotonic config could end up with NO next
# trigger until a reboot.
[Unit]
Description=Run the Trusty Squire registry verify pass daily

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
# Spread the start so it never lines up exactly with other timers.
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
