# tools/harvester/systemd/harvester-digest.service
[Unit]
Description=Trusty Squire skill harvester — daily Telegram digest
Documentation=https://github.com/Trusty-Squire/trusty-squire/blob/main/tools/harvester/README.md

[Service]
Type=oneshot
WorkingDirectory=/home/lunchbox/proj-ts/tools/harvester
ExecStart=/usr/bin/node /home/lunchbox/proj-ts/tools/harvester/daily-digest.mjs

# Needs TELEGRAM_BOT_TOKEN. Without it the script prints to stderr
# (and exits 1) so the systemd OnFailure= hook (future Phase 4 work)
# can alert via an alternative channel.
EnvironmentFile=-%h/.config/trusty-squire/harvester.env

StandardOutput=journal
StandardError=journal
