# tools/harvester/systemd/harvester-daily-open.timer
[Unit]
Description=Trusty Squire skill harvester — open daily summary at 00:05 UTC

[Timer]
OnCalendar=*-*-* 00:05:00 UTC
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
